# store/config/agentd/service/housing-collector.timer
[Unit]
Description=Housing Data Collection Timer
Requires=housing-collector.service

[Timer]
# Run twice daily: 9 AM and 6 PM
OnCalendar=*-*-* 09:00:00
OnCalendar=*-*-* 18:00:00
# Also run 20 minutes after boot
OnBootSec=20min
# If missed, run immediately
Persistent=true

[Install]
WantedBy=timers.target
