# .mark/ai-agents/Codriver/etc/systemd/lead-scraper-daily.timer
[Unit]
Description=Daily Lead Scraping Timer
Requires=lead-scraper-daily.service

[Timer]
# Run every day at 9 AM
OnCalendar=daily
OnCalendar=*-*-* 09:00:00
Persistent=true

[Install]
WantedBy=timers.target
