# .mark/ai-agents/Marketeer/etc/systemd/marketeer.service
[Unit]
Description=Marketeer Edge Router
After=network.target codriver.service
Wants=network.target codriver.service

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/ws/.mark/ai-agents/Marketeer
EnvironmentFile=-/home/admin/ws/.mark/ai-agents/Marketeer/.env
ExecStart=/home/admin/.cargo/bin/cargo run --bin marketeer
Restart=always
RestartSec=10
StandardOutput=append:/var/log/marketeer.log
StandardError=append:/var/log/marketeer-error.log

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target