# freightdev/MarkOS/lib/Codriver/etc/systemd/codriver.service
[Unit]
Description=CoDriver API Gateway & AI Orchestrator
After=network.target
Wants=network.target

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/ws/.mark/ai-agents/Codriver
EnvironmentFile=-/home/admin/ws/.mark/ai-agents/Codriver/.env
ExecStart=/home/admin/.cargo/bin/cargo run --bin codriver
Restart=always
RestartSec=10
StandardOutput=append:/var/log/codriver.log
StandardError=append:/var/log/codriver-error.log

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target